# crontab
0 * * * * /usr/local/bin/monitor.sh
0 1 4,20 * * /usr/local/bin/backup.sh
0 1 1 * * /usr/local/bin/cleanup.sh